# xfce4-session.service
[Unit]
BindsTo=graphical-session.target

[Service]
ExecStart=/usr/bin/startxfce4
